# files/service_files/disable_thp.timer
[Unit]
Description=Ensure that THP is disabled at boot

[Timer]
OnBootSec=0

[Install]
WantedBy=timers.target
