# files/service_files/disable_thp.service
[Unit]
Description=Disables THP on boot

[Service]
Type=oneshot
ExecStart=/bin/hugeadm --thp-never

[Install]
WantedBy=default.target
